# pi-monitor/streamdesk-monitor.service
# StreamDesk Monitor — systemd service for Raspberry Pi
#
# Setup:
#   1. Install: pip install playwright && playwright install chromium
#   2. Edit ExecStart below: set SERVER_URL and DEVICE_KEY
#   3. Copy: sudo cp streamdesk-monitor.service /etc/systemd/system/
#   4. Enable: sudo systemctl enable streamdesk-monitor
#   5. Start:  sudo systemctl start streamdesk-monitor
#
# Logs:
#   journalctl -u streamdesk-monitor -f

[Unit]
Description=StreamDesk Monitor
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
User=pi
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
ExecStart=/usr/bin/python3 /home/pi/StreamDesk-Agent/pi-monitor/monitor.py http://SERVER_IP:8000 MY_DEVICE_KEY
Restart=on-failure
RestartSec=10

[Install]
WantedBy=graphical-session.target
